# deployment/remote-agent/services/openclaw-service/openclaw-service.service
[Unit]
Description=OpenClaw AI Agent Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/openclaw-service
ExecStart=/usr/bin/node /opt/openclaw-service/src/index.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=openclaw-service

[Install]
WantedBy=multi-user.target
